# config/services/cine4_data.service
[Service]
ExecStart=/usr/bin/php /home/cine/cine4/src/dataServer.php
